# data/org.gnome.Papers.service.in
[Unit]
Description=Papers document viewer
Documentation=man:papers(1)

[Service]
Type=dbus
BusName=org.gnome.papers.Daemon
ExecStart=@libexecdir@/papersd
